# chkphy
#! /bin/bash

### BEGIN INIT INFO
# Provides:          chkphy
# Required-Start:    
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Check BBB phy on boot and repower if PHY failed
# Description:       This file should be used to construct scripts to be
#                    placed in /etc/init.d.  This example start a
#                    single forking daemon capable of writing a pid
#                    file.  To get other behavoirs, implemend
#                    do_start(), do_stop() or other functions to
#                    override the defaults in /lib/init/init-d-script.
### END INIT INFO

# The following part always gets executed.


# first we sleep.
# this gives plently of time for linux to mess with the RTC 
# (where does it do that set?) before we touch it. Otherwise it 
# occasionally updates it right between when we are doing our sets.

# This also gives you a chance to abort in case something goes wrong
# and avoid a infinite loop

sleep 30

if [[ $(phyreg test 18 13) == "1" ]]; then 

	sleep 1

	NOW=$(bbbrtc now)
	echo "chkphy:Rebooting NOW=$NOW" >/dev/kmsg
	echo "chkphy:Rebooting NOW=$NOW" >/dev/console
	echo "Rebooting" >>/var/tmp/chkphy.log
	# make sure that log message actually gets written before we pull the plug
	sync
	i2cset -f -y 0 0x24 0x0a 0x00
	bbbrtc now 130
	bbbrtc wake 110
	bbbrtc sleep 102
	bbbrtc now 100
	while true; 
	  do echo waiting to reboot
	  sleep 10
	done
fi

echo "chkphy:eth0 good" >/dev/kmsg
echo "chkphy:eth0 good" >/dev/console
